# makesite/templates/uwsgi/deploy/supervisor.conf
[program:{{ project }}.{{ branch }}]
command=/usr/bin/uwsgi
  --socket {{ deploy_dir }}/uwsgi.sock
  --home {{ deploy_dir }}/.virtualenv/
  --pythonpath {{  project_sourcedir  }}
  --module wsgi
  --processes 1
  --master
  --harakiri 120
  --max-requests 5000
user={{ site_user }}
environment=DJANGO_SETTINGS_MODULE="settings{{ '.' + mode if mode else '' }}"
stdout_logfile={{ deploy_dir }}/logs/uwsgi.log
stderr_logfile={{ deploy_dir }}/logs/uwsgi_err.log
directory={{ project_sourcedir }}
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
