# makesite/templates/uwsgi/deploy/supervisor.conf
[program:{{ project }}.{{ branch }}]
command=uwsgi -x {{ uwsgi_source_confpath }}
user={{ site_user }}
environment=DJANGO_SETTINGS_MODULE="settings{{ '.' + mode if mode else '' }}"
stdout_logfile={{ deploy_dir }}/logs/uwsgi.log
stderr_logfile={{ deploy_dir }}/logs/uwsgi_err.log
directory={{ project_sourcedir }}
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
